# systemd/codeframe-staging.service
[Unit]
Description=CodeFRAME Staging Server
Documentation=https://github.com/frankbria/codeframe
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=frankbria
WorkingDirectory=/home/frankbria/projects/codeframe
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/frankbria/projects/codeframe/.env.staging
ExecStart=/usr/bin/pm2 start /home/frankbria/projects/codeframe/ecosystem.staging.config.js
ExecStop=/usr/bin/pm2 stop all
ExecStop=/bin/sleep 5
ExecStop=/usr/bin/pm2 delete all
ExecReload=/usr/bin/pm2 restart all
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
